# .docker/crontab
* * * * * cd /var/www && php artisan schedule:run 2>&1 | grep -v "No scheduled commands are ready to run."